# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ─── PHP-FPM ──────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=1
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

; ─── Nginx ────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=2
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

; ─── Queue Worker ─────────────────────────────────────────────────
[program:queue-worker]
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
priority=10
user=www-data
stderr_logfile=/var/log/queue.err.log
stdout_logfile=/var/log/queue.out.log

; ─── Scheduler ────────────────────────────────────────────────────
[program:scheduler]
command=sh -c "while true; do php /var/www/html/artisan schedule:run; sleep 60; done"
autostart=true
autorestart=true
priority=10
user=www-data
stderr_logfile=/var/log/scheduler.err.log
stdout_logfile=/var/log/scheduler.out.log
